# root/etc/services.d/tor/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid abc /entrypoint.sh "${RUN_OPTS}"

